# standalone/crontab
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/bin

*/15 * * * * echo "check syntax fixtures"
